# services/postgresql.service.d/override.conf
[Service]
Environment=PGDATA=/mnt/linux_data/Database/postgresql/data
